# repo_sync.cron
0 0 * * * /reposado/code/repo_sync --log /reposado/metadata/repo_sync.log
